# dockerfiles/icat/supervisord.conf
[supervisord]
nodaemon=true

[program:ndslabs-irods-federate]
command=/opt/ndslabs/ndslabs-irods-federate -zone RODS_ZONE -host HOSTNAME
